# hometv-master/hometv-master-server/src/deb/control/postinst
#!/bin/sh
update-rc.d hometv-master-server defaults
/etc/init.d/hometv-master-server start
