# hometv-master/hometv-master-server/src/deb/control/prerm
#!/bin/sh
/etc/init.d/hometv-master-server stop
update-rc.d hometv-master-server remove
